# contrib/systemd/gmail-relay.path
[Path]
PathExistsGlob=/var/nullmailer/queue/*

[Install]
WantedBy=multi-user.target
